# workflow/com.promptlab.synthesizer.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.promptlab.synthesizer</string>

    <key>ProgramArguments</key>
    <array>
        <string>__PYTHON3__</string>
        <string>__REPO_DIR__/synthesizer.py</string>
        <string>--all</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>2</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>__REPO_DIR__/synthesizer.log</string>

    <key>StandardErrorPath</key>
    <string>__REPO_DIR__/synthesizer.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
